# crontab
#Respaldar la carpeta de origen a /backup diariamente a las 2 a.m

0 2 * * * rsync -a  /home/nico18/so /home/nico18/so/Backup
